# root/etc/services.d/monit/run
#!/usr/bin/with-contenv bash

# Run monit
/usr/bin/monit -I -c /etc/monitrc

