# root/etc/s6-overlay/s6-rc.d/monit/run
#!/usr/bin/with-contenv bash

# Run monit
exec /usr/bin/monit -I -c /etc/monitrc
